# iblock.rc
#!/bin/ksh

daemon="/usr/local/sbin/iblock"
daemon_user="_iblock"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
